# Docker-compose/django_nginx_redis/django/supervisord.conf
[program:server]
command=uwsgi /opt/server/uwsgi.ini
directory=/opt/server
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true